# crontab
0 */2 * * * /home/marc/tool/websupportsk-ddns/run-sync
